# gogs-codex-arm-env-20260419-100157/project/GOGS/deploy/systemd/grab-system.service
[Unit]
Description=Grab Inventory Qt Backend Service
After=network-online.target mysql.service mediamtx.service
Wants=network-online.target

[Service]
Type=simple
User=__RUN_USER__
Group=__RUN_USER__
WorkingDirectory=/userdata/GOGS/backend
ExecStart=/bin/bash /userdata/GOGS/backend/start.sh
Restart=always
RestartSec=5
LimitNOFILE=65535

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictRealtime=true
MemoryDenyWriteExecute=true
LockPersonality=true

# Read-write paths (application data directories)
ReadWritePaths=/userdata/GOGS/backend/config
ReadWritePaths=/userdata/GOGS/backend/logs
ReadWritePaths=/userdata/GOGS/backend/cache
ReadWritePaths=/userdata/GOGS/backend/tmp
ReadWritePaths=/userdata/GOGS/backend/data
ReadWritePaths=/userdata/GOGS/backend/maintenance
ReadWritePaths=/tmp

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=grab-system

[Install]
WantedBy=multi-user.target
